# hp-printer-keepalive.service
[Unit]
Description=HP printer keepalive
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/hp-printer-keepalive.sh
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
